# conf/appmon.service
[Unit]
Description=appmon
After=syslog.target

[Service]
ExecStart=/opt/appmon/appmon.jar --server.port=34567

[Install]
WantedBy=multi-user.target
